# sguclient-OpenWrt-LuCI/sguclient/luci_files/root/etc/init.d/sguclient
#!/bin/sh /etc/rc.common
START=99

run_sguclient()
{
	local enable
	config_get_bool enable $1 enable

	if [ $enable ]; then
		local autoreconnect
		local username
		local password
		local isptype
		local ifname
		local wanip

		
		config_get_bool autoreconnect $1 autoreconnect
		config_get_bool noheartbeat $1 noheartbeat
		config_get username $1 username
		config_get password $1 password
		config_get isptype $1 isptype
		config_get ifname $1 ifname
		config_get wanip $1 wanip

		
		process=`pgrep sguclient`

		if [ -n "$process" ];then
			echo "sguclient is running, ignored..."
		else
			echo "sguclient not running, trying to start again..."
			
			local argconfig=""

			#----To start sguclient
			if [ $autoreconnect ]; then
				argconfig="--auto"
			fi

			if [ $noheartbeat ]; then
				argconfig="$argconfig --noheartbeat"
			fi
			
			rm -rf /etc/crontabs/root
			
			# if active then create guard scripts
			sguclient -u $username -p $password -i $isptype --device $ifname --ip $wanip $argconfig &
			echo "*/1 * * * * sguclient -u $username -p $password -i $isptype --device $ifname --ip $wanip $argconfig" >> /etc/crontabs/root
			
			# TODO:auto config network
			local extranetUsername
			local extranetPassword
			config_get extranetUsername $1 extranetUsername
			config_get extranetPassword $1 extranetPassword
			# wan interface
			uci set network.wan.ifname=$ifname
			uci set network.wan._orig_ifname=$ifname
			uci set network.wan._orig_bridge='false'
			uci set network.wan.proto='static'
			uci set network.wan.ipaddr=$wanip
			uci set network.wan.netmask='255.255.255.0'
			uci set network.wan.gateway=$(echo $wanip|cut -d . -f1-3).254

			hasRoute=$(uci get network.@route[0])
			if [ "$hasRoute" != "route" ];then
				# route1
				uci add network route
				uci set network.@route[-1].interface='wan'
				uci set network.@route[-1].target='210.38.0.0'
				uci set network.@route[-1].netmask='255.255.255.0'
				uci set network.@route[-1].gateway=$(echo $wanip|cut -d . -f1-3).254

				# route2
				uci add network route
				uci set network.@route[-1].interface='wan'
				uci set network.@route[-1].target='192.168.0.0'
				uci set network.@route[-1].netmask='255.255.255.0'
				uci set network.@route[-1].gateway=$(echo $wanip|cut -d . -f1-3).254
			fi

			# extranet
			uci set network.wan6.ifname=$ifname
			uci set network.wan6._orig_ifname=$ifname
			uci set network.wan6._orig_bridge='false'
			uci set network.wan6.proto='pppoe'
			uci set network.wan6.ipv6='auto'
			uci set network.wan6.username=$extranetUsername
			uci set network.wan6.password=$extranetPassword

			uci commit network
			/etc/init.d/network restart

			echo "SGUClient has started."
			#----

		fi
	else
		rm -rf /etc/crontabs/root
	fi
}

start()
{
	config_load sguclient
	config_foreach run_sguclient login
}

stop()
{
	killall sguclient
	rm -rf /etc/crontabs/root
	echo "SGUClient has stoped."
}
